# tools/webdash.service
[Unit]
Description=NetDash Web Dashboard
After=network.target

[Service]
Type=simple
User=clawdbot
ExecStart=/usr/bin/python3 /home/clawdbot/aispace/tools/webdash.py 8088
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
